# stage2.custom/01-tweaks/files/rc.local
## Custom: Everything below this line is auto-edited by setup scripts

/opt/docker/bin/docker_run.sh | tee /tmp/docker_run.log
/opt/iptables/iptables_rules.sh
